# systemd/lora-gateway-bridge.service.d/lora-gateway-bridge.conf
[Service]
Environment="SETTING_A=value_a"
Environment="SETTING_B=value_b"
